# chromiarchos.conf
# Copyright (c) 2011 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

start on starting system-services
stop on stopping system-services

exec /usr/local/bin/chromiarchos-run

post-start script
  # Disable screen blanking and powering down for dev & factory.
  # TODO(tbroch) If X is running and DPMS is enabled it still blanks VT2.
  # Investigate possibility of tying these different code paths together.
  setterm -cursor on -blank 0 -powersave off -powerdown 0 >/dev/tty3 || true
end script
